# deploy/systemd/loom.service
[Unit]
Description=Loom Production Stack
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/loom/compose
ExecStart=/usr/bin/docker compose -f /opt/loom/compose/docker-compose.production.yml --env-file /opt/loom/env/.env.production up -d --remove-orphans
ExecStop=/usr/bin/docker compose -f /opt/loom/compose/docker-compose.production.yml --env-file /opt/loom/env/.env.production down
ExecReload=/usr/bin/docker compose -f /opt/loom/compose/docker-compose.production.yml --env-file /opt/loom/env/.env.production up -d --remove-orphans
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
