# deploy/systemd/loom.service
[Unit]
Description=Loom Production Stack
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/lingfeng/loom
ExecStart=/usr/bin/docker compose -f /home/lingfeng/loom/docker-compose.yml --env-file /home/lingfeng/loom/.env up -d --no-build --remove-orphans
ExecStop=/usr/bin/docker compose -f /home/lingfeng/loom/docker-compose.yml --env-file /home/lingfeng/loom/.env down
ExecReload=/usr/bin/docker compose -f /home/lingfeng/loom/docker-compose.yml --env-file /home/lingfeng/loom/.env up -d --no-build --remove-orphans
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
