# deploy/systemd/template.service
[Unit]
Description=Template Production Stack
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/template/compose
ExecStart=/usr/bin/docker compose -f /opt/template/compose/docker-compose.production.yml --env-file /opt/template/env/.env.production up -d --remove-orphans
ExecStop=/usr/bin/docker compose -f /opt/template/compose/docker-compose.production.yml --env-file /opt/template/env/.env.production down
ExecReload=/usr/bin/docker compose -f /opt/template/compose/docker-compose.production.yml --env-file /opt/template/env/.env.production up -d --remove-orphans
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
